# deployment/etc/init.d/pdt.conf
start on filesystem and net-device-up IFACE=lo
exec /usr/lib/pdt/bin/circusd /etc/pdt/circus.ini
