# deployment/etc/init/pdt.conf
start on filesystem and net-device-up IFACE=lo
chdir /usr/lib/pdt
exec python3 /usr/lib/pdt/bin/circusd /etc/pdt/circus.ini
